# sickbeard/files/supervisor.conf
[supervisord]
nodaemon=true

[program:app]
priority=10
directory=/opt/sickbeard
command=python /opt/sickbeard/SickBeard.py --datadir=/config
user=nobody
autorestart=true
autostart=true
